# config/files/usr/etc/xdg/autostart/1Password-silverflow.desktop
[Desktop Entry]
Type=Application
Name=1Password
Comment=1Password silent autostart
Icon=1password
Exec=/opt/1Password/1password --silent 